# chat.conf
[program:chat]
command=/usr/local/bin/gunicorn wsgi -c gunicorn.config.py
directory=/var/www/chat
autostart=true
autorestart=true